# rootfs/etc/crontabs/crontab
# Run Vault Backups every hour at X:10
10 * * * * flock -xn /tmp/backup.lock /usr/local/bin/vault-backup.sh
