# crontabs/crontab
# Run once every hour
@hourly /usr/local/bin/vault-migrate-backup.sh
