# crontabs/crontab
# Run once every hour
@hourly /usr/local/bin/vault-migrate-gcs-file.sh
